# infrastructure/alpha-heartbeat.timer
[Unit]
Description=Alpha Heartbeat Timer - Periodic activation for Project Beta

[Timer]
# Midnight to 4 AM, every 15 minutes (17 heartbeats)
# First night: December 8-9, 2025
OnCalendar=*-*-* 00:00,00:15,00:30,00:45,01:00,01:15,01:30,01:45,02:00,02:15,02:30,02:45,03:00,03:15,03:30,03:45,04:00:00

# Don't run missed heartbeats if the Pi was asleep
Persistent=false

# Add some randomness to avoid exact-second firing
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
