# infrastructure/alpha-heartbeat.service
[Unit]
Description=Alpha Heartbeat - Wake Alpha for time alone
After=network.target

[Service]
Type=oneshot
User=jefferyharrell
WorkingDirectory=/home/jefferyharrell/Projects/Alpha-Home
ExecStart=/home/jefferyharrell/Projects/Alpha-Home/wake-alpha.sh
StandardOutput=append:/home/jefferyharrell/Projects/Alpha-Home/logs/heartbeat.log
StandardError=append:/home/jefferyharrell/Projects/Alpha-Home/logs/heartbeat.log

# Don't let a single failed heartbeat stop future ones
Restart=no
